# hardware/pi-node/systemd/openpango.service
[Unit]
Description=OpenPango Agent Node
Documentation=https://github.com/openpango/openpango-skills
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=openpango
Group=openpango
WorkingDirectory=/opt/openpango
Environment="NODE_ENV=production"
Environment="OLLAMA_HOST=http://localhost:11434"
ExecStart=/usr/bin/openpango start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
